# noesis.service
[Unit]
Description=Noesis — Decentralized Cognitive Architecture
Documentation=https://github.com/hitensaxena/Noesis
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=hitensaxena
ExecStart=/Users/hitensaxena/.cargo/bin/cargo run --release --bin noesis -- start --rest --port 8647 --storage postgres --database-url "host=127.0.0.1 port=5433 dbname=noesis user=noesis password=noesis" --redis-url "redis://127.0.0.1:6380"
WorkingDirectory=/Users/hitensaxena/noesis
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=noesis=info

[Install]
WantedBy=multi-user.target
